# kinome-database/website_files/etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)
[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:kinomeDBAPI]
command=uwsgi -s /tmp/uwsgi.sock -w kinomeDBAPI:app
user=nginx
directory=/var/www/kinomeDBAPI
autostart=true
autorestart=true
stdout_logfile=/var/log/uwsgi/kinomeDBAPI.log
redirect_stderr=true
stopsignal=QUIT

;[program:testAPI]
;command=uwsgi -s /tmp/uwsgi.sock -w myapp:app
;user=ec2-user
;directory=/var/www/testAPI
;autostart=true
;autorestart=true
;stdout_logfile=/var/log/uwsgi/testAPI.log
;redirect_stderr=true
;stopsignal=QUIT
